# com.user.fixwallpaper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.fixwallpaper</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/REPLACE_WITH_USERNAME/bin/wp_overlay</string>
        <string>REPLACE_WITH_VIDEOS_DIR</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/Library/Logs/com.user.fixwallpaper.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/REPLACE_WITH_USERNAME/Library/Logs/com.user.fixwallpaper.err.log</string>
</dict>
</plist>
